# templates/freerds-manager.service.erb
[Unit]
Description=FreeRDS manager
After=syslog.target network.target

[Service]
ExecStart=/opt/FreeRDS/bin/freerds-manager --no-daemon

[Install]
WantedBy=multi-user.target
