# templates/freerds-server.service.erb
[Unit]
Description=FreeRDS server
After=syslog.target network.target

[Service]
ExecStart=/opt/FreeRDS/bin/freerds-server --no-daemon

[Install]
WantedBy=multi-user.target
